# layout/Library/LaunchDaemons/com.miwix.dynamicpaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miwix.dynamicpaper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/dynamicPaper</string>
    </array>
    <key>StartInterval</key>
    <integer>86400</integer>
    <key>StandardOutPath</key>
    <string>/var/mobile/dynamicPaper.log</string>
    <key>StandardErrorPath</key>
    <string>/var/mobile/dynamicPaper.log</string>
    <key>Debug</key>
    <true/>
    <key>SoftResourceLimits</key>
    <dict>
        <key>Core</key>
        <integer>9223372036854775807</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>Core</key>
        <integer>9223372036854775807</integer>
    </dict>
	<key>UserName</key>
	<string>mobile</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
	<true/>
</dict>
</plist>
