# layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.miwix.dynamicpaper.plist
